# crontab
0 2 * * * {{MFSERV_HOME}}/bin/cronwrap.sh --lock --log-capture-to={{MFMODULE_RUNTIME_HOME}}/log/changeslog.log -- plugin_wrapper generate_changelogs changelogs.sh
